# etc/cron.d/certbot-weekly
2 5 * * 0 <username> cd /home/<username>/matrix-server && docker compose up certbot
